# gate.service
[Unit]
Description=Gate service
After=network.target

[Service]
ExecStart=/usr/bin/gate
Restart=always
User=pi

[Install]
WantedBy=multi-user.target